# services/powertop.service
[Unit]
Description=Powertop tunings

[Service]
Type=oneshot
ExecStart=/home/minhduc0711/.local/bin/powertop.sh

[Install]
WantedBy=multi-user.target
